# gunicorn/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=root
childlogdir=/tmp

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:gunicorn]
command=/usr/local/bin/gunicorn -c /home/irocha/git/python-labs/gunicorn/gunicorn.conf.py --log-config /home/irocha/git/python-labs/gunicorn/logging-console.conf test:app
directory=/home/irocha/git/python-labs/gunicorn
autostart=true
autorestart=true
redirect_stderr=True
stdout_logfile=/tmp/test.log
stdout_logfile_maxbytes=10MB
