# supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=root
childlogdir=/tmp/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

# supervisorctl start|stop web:
[group:web]
programs=web-7000,web-7001

[program:web-7000]
command=/home/irocha/git/python-labs/tornado/web/Web.py --logging=debug --port=7000
redirect_stderr=true
stdout_logfile=/tmp/web-7000.log
stdout_logfile_maxbytes=10MB

[program:web-7001]
command=/home/irocha/git/python-labs/tornado/web/Web.py --logging=debug --port=7001
redirect_stderr=true
stdout_logfile=/tmp/web-7001.log
stdout_logfile_maxbytes=10MB
